# bin/backup/db2psc.crontab
MAILTO=""
*/15 * * * * /home/db2psc/bin/liinfo
